# systemd/virtual-revelation-world-pack.timer
[Unit]
Description=Virtual Revelation World Pack Timer

[Timer]
OnCalendar=hourly
Persistent=false
AccuracySec=1m

[Install]
WantedBy=timers.target
